# system_files/fwload/usr/lib/systemd/system/ModemManager.service.d/wait-for-fw.conf
[Unit]
After=sdm845-fwload.service
Requires=sdm845-fwload.service
